# motd.service
[Unit]
Description=Update MOTD

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/update-motd
